# docs/benchmarks/mkosi-v25/mkosi.extra/etc/systemd/system/serial-getty@.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noclear --autologin root --keep-baud %I $TERM

[Install]
Alias=getty.target.wants/serial-getty@ttyS0.service
